# common/.config/autostart/background_image_downloader.desktop
[Desktop Entry]
Name=BACKGROUND IMAGE DOWNLOADER
Exec=bash download_background_image.sh
Terminal=false
Type=Application
