# common/.config/autostart/dms_upload_runner.desktop
[Desktop Entry]
Name=DMS UPLOAD RUNNER
Exec=bash dms_upload.sh
Terminal=false
Type=Application
