# common/.config/autostart/scrcpy_runner.desktop
[Desktop Entry]
Name=SCRCPY RUNNER
Exec=bash run_scrcpy.sh
Terminal=false
Type=Application
